# recipes-crypto/garbage-collector/files/garbage-collector.timer
[Unit]
Description=Garbage Collector Timer
After=fcc-crypto.service

[Timer]
OnUnitActiveSec=30 days

[Install]
WantedBy=timers.target
