# recipes-crypto/garbage-collector/files/garbage-collector.service
[Unit]
Description=Garbage Collector Service

[Service]
Type=oneshot
ExecStart=/bin/garbage-collector.sh

[Install]
WantedBy=timers.target
